# packages/mediacenter/xbmc/system.d/xbmc.target
[Unit]
Description=XBMC Media Center
Conflicts=rescue.service rescue.target
Requires=graphical.target
After=graphical.target
AllowIsolate=yes
