# packages/network/openssh/system.d/sshd-keygen.service
[Unit]
Description=SSH server keys generation.
After=syslog.target
Before=sshd.service

[Service]
Type=oneshot
ExecStart=/usr/lib/openssh/sshd-keygen
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
