# packages/sysutils/busybox/system.d/debug-shell.service
[Unit]
Description=Debug Shell
ConditionKernelCommandLine=debugging

[Service]
Environment=HOME=/storage TERM=linux
WorkingDirectory=/storage
ExecStart=/sbin/getty -n 38400 tty3
Restart=always
RestartSec=0
#UtmpIdentifier=$TTY
#TTYPath=/dev/tty$TTY
#TTYReset=yes
#TTYVHangup=yes
#TTYVTDisallocate=yes
KillMode=process

# Some login implementations ignore SIGTERM, so we send SIGHUP
# instead, to ensure that login terminates cleanly.
KillSignal=SIGHUP

[Install]
WantedBy=multi-user.target
